# services/perconadb-tokudb/supervisord.conf
[program:perconadb]
command=mysqld_safe
